# deploy/unifiedops-health-check.service
; ============================================================================
; UnifiedOps health-check service — health_check.py
;
; Responsibilities:
;   - Polls all 3 per-server heartbeat InfluxDB buckets every 10 s.
;   - Tracks per-listener UP / DOWN / INFRA_DOWN state with stable
;     event_keys so each outage raises a modal exactly once.
;   - Pushes the snapshot + raised events to React via WebSocket
;     at /ws/listener-health on the same FastAPI process.
;
;   * Runs *next to* server.py on the UI VM; both share the multithreaded
;     httpx pool. Splitting them is a code-level concern only.
;
; Install:
;   sudo cp deploy/unifiedops-health-check.service /etc/systemd/system/
;   sudo systemctl daemon-reload
;   sudo systemctl enable --now unifiedops-health-check
;
; ============================================================================
[Unit]
Description=UnifiedOps v2 health-check (listener heartbeat WebSocket)
After=network-online.target unifiedops-ui-server.service
Wants=network-online.target
PartOf=unifiedops-ui-server.service

[Service]
Type=simple
User=unifiedops
Group=unifiedops
WorkingDirectory=/opt/unifiedops
EnvironmentFile=/etc/unifiedops/ui.env
ExecStart=/opt/unifiedops/.venv/bin/python /opt/unifiedops/server/health_check.py
Restart=on-failure
RestartSec=4
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/log/unifiedops

[Install]
WantedBy=multi-user.target
